# host-temp-sender.service
[Unit]
Description=Host Temperature Sender via virtio-serial
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/host-temp-sender.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
